# deploy/shillbot-close-11pm.service
[Unit]
Description=Shooter ShillBot Window Close (11:00 PM CT)
After=network.target

[Service]
Type=oneshot
User=%i
WorkingDirectory=/home/%i/shillbot
EnvironmentFile=-/home/%i/shillbot/.env
ExecStart=/home/%i/shillbot/.venv/bin/python -m shillbot close-once
StandardOutput=journal
StandardError=journal
